# etc/conf/alpine/mos-guest/includes.chroot/etc/local.d/0250-networking.start
#!/bin/sh

echo "nameserver 10.0.4.2" > /etc/resolv.conf
