# etc/conf/alpine/mos-fw/includes.chroot/etc/local.d/0050-devconf.start
#/bin/sh

mkdir /dev/pts
mkdir /dev/shm

mount -t devpts -o gid=4,mode=620 none /dev/pts
mount -t tmpfs none /dev/shm
