# conf/alpine/mos-guest/includes.chroot/etc/local.d/0225-ssh-keys-conf.start
#!/bin/sh

mkdir -p /home/user/.ssh/
cp /etc/ssh/authorized_keys /home/user/.ssh/authorized_keys

chown -R /home/user/.ssh user:user
